# templates/default/upstart_archipel.conf.erb
# Archipel - The real time orchestrator
#
# Archipel-agent

description "Archipel Agent"

start on started libvirt-bin 
stop on stopping libvirt-bin

respawn
respawn limit 5 30

exec runarchipel --config=/etc/archipel/archipel.conf -n
